# systemd/attic.service
[Unit]
Description=Attic backup

[Service]
Type=oneshot
ExecStart=/home/xconstruct/code/conf/dotfiles/bin/attic-home.sh
